# nonnix/luci-theme-we/root/etc/uci-defaults/luci-theme-we
#!/bin/sh
uci batch <<-EOF
	set luci.themes.we=/luci-static/we
	set luci.main.mediaurlbase=/luci-static/we
	commit luci
EOF
exit 0
