# nonnix/luci-theme-xwifi/root/etc/uci-defaults/luci-theme-xwifi
#!/bin/sh
uci batch <<-EOF
	set luci.themes.XWiFi=/luci-static/xwifi
	set luci.main.mediaurlbase=/luci-static/xwifi
	commit luci
EOF
exit 0
