# nonnix/luci-theme-simple/root/etc/uci-defaults/30_luci-theme-simple
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Simple=/luci-static/simple
	set luci.main.mediaurlbase=/luci-static/simple
	commit luci
EOF
exit 0
