# nonnix/luci-theme-vi/root/etc/uci-defaults/30_luci-theme-vi
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci batch <<-EOF
		set luci.themes.vi=/luci-static/vi
		set luci.main.mediaurlbase=/luci-static/vi
		commit luci
	EOF
fi

exit 0
